# modules/genomics_cromwell/scripts/cromwell.service
[Unit]
Description=Cromwell Server
After=network.target

[Service]
User=root
Group=root
Restart=always
TimeoutStopSec=10
RestartSec=5
WorkingDirectory=/opt/cromwell
ExecStart=/usr/bin/java -Xmx10G -Dconfig.file=cromwell.conf -jar cromwell.jar server

[Install]
WantedBy=multi-user.target